# system-monitor/launchd/com.trey.system-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.trey.system-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/trey/Desktop/Apps/system-monitor/dist/index.js</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/trey/Desktop/Apps/system-monitor/logs/daemon-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/trey/Desktop/Apps/system-monitor/logs/daemon-stderr.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/trey/Desktop/Apps/system-monitor</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>

  <key>LowPriorityBackgroundIO</key>
  <true/>

  <key>Nice</key>
  <integer>10</integer>
</dict>
</plist>
